# examples/proxmcp.service
[Unit]
Description=proxmcp MCP server
After=network.target

[Service]
Type=simple
User=root
Environment=PROXMCP_AUDIT_LOG=/var/log/proxmcp-audit.log
Environment=PROXMCP_ALLOW_BREAK_GLASS=1
ExecStart=proxmcp-server
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
